# deploy/bingo-bot.service
[Unit]
Description=Bingo Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bingo
ExecStart=/home/ubuntu/bingo/.venv/bin/python /home/ubuntu/bingo/bot.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
