# deploy/bingo-bot.service
[Unit]
Description=Bingo Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opc
WorkingDirectory=/home/opc/bingo
ExecStart=/home/opc/bingo/.venv/bin/python /home/opc/bingo/bot.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
